# deploy/polaris.service
[Unit]
Description=Polaris — IP Management Tool
Documentation=https://github.com/davidmoore-rogers/polaris
# PGDG packages name the unit `postgresql-15.service`. Adjust the version
# suffix to match your install (`postgresql-16.service`, etc.). If you're on
# the RHEL AppStream module instead of PGDG, change both lines below to
# reference `postgresql.service` — see docs/INSTALL.md for the full
# migration path from AppStream to PGDG.
After=network.target postgresql-15.service
Requires=postgresql-15.service

[Service]
Type=simple
User=polaris
Group=polaris
WorkingDirectory=/opt/polaris
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/polaris/.env
Restart=on-failure
RestartSec=5

# Security hardening
AmbientCapabilities=CAP_NET_BIND_SERVICE
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/polaris
PrivateTmp=true

# Logging — stdout/stderr go to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polaris

[Install]
WantedBy=multi-user.target
